# deploy/supervisord.conf
; klub-lotto pod process supervisor.
;
; This file is loaded by Debian's default /etc/supervisor/supervisord.conf
; via the [include] files=/etc/supervisor/conf.d/*.conf directive. We
; only declare [program:*] sections here — a second [supervisord] section
; would conflict with the one in the main config. To run in the foreground
; (required in a container) entrypoint.sh passes `-n` to supervisord.
;
; Order matters only because Xvfb needs to be up before Chromium can be
; launched by agent-browser. supervisord starts processes in priority
; order; lower number = earlier.

[program:xvfb]
priority=10
command=/usr/bin/Xvfb :99 -screen 0 1600x1200x24 -ac -nolisten tcp
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

; Tiny window manager so Chromium's window can be resized and focused.
[program:fluxbox]
priority=20
environment=DISPLAY=:99
command=/usr/bin/fluxbox -rc /root/.fluxbox/init -no-slit -no-toolbar 2>/dev/null || /usr/bin/fluxbox -rc /root/.fluxbox/init
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

; Keep the virtual desktop clean. Fluxbox sometimes tries fbsetbg on startup
; (especially on first run or after PVC restore), which pops an xmessage dialog.
; We run a dedicated script that forces a solid dark background and kills
; stray dialogs in a loop. Using an external script avoids horrific quoting
; problems inside the supervisord .ini file.
[program:desktop-cleanup]
priority=25
environment=DISPLAY=:99
command=/usr/local/bin/desktop-cleanup.sh
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:x11vnc]
priority=30
; -bg would daemonize, which we don't want under supervisord (we need
; x11vnc to stay in the foreground so supervisord can track it). Just
; omit the flag — foreground is the default.
command=/usr/bin/x11vnc -display :99 -forever -shared -nopw -rfbport 5900 -noxdamage -xkb -noxrecord -noxfixes -wait 5 -defer 5 -nodpms -nossl
; NOTE: -quiet intentionally omitted (was added for DPMS/charset suppression per follow-up 3) to preserve WS/novnc error messages in logs for debugging VNC issues the hardening was meant to address. Other flags + fluxbox init cover warnings.
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

; websockify bundles noVNC's index.html at /. We serve it on 6080 and the
; web server proxies /vnc/* through to it.
[program:novnc]
priority=40
command=/usr/bin/websockify --web /usr/share/novnc 6080 127.0.0.1:5900
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:web]
priority=50
environment=KLUBLOTTO_HEADED="true",DISPLAY=":99"
command=/usr/local/bin/klub-lotto-web
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
